# dbipc/rootfs/etc/services.d/rabbitmq/finish
#!/bin/execlineb -P

define pid /var/run/rabbitmq/pid

foreground {
    if { test -e ${pid} }
    rm ${pid}
}
